# demo-04-full/proxy/root/etc/services.d/rrproxy/run
#!/usr/bin/with-contenv sh

# start up the proxy
exec /usr/local/bin/rrproxy --config /root/config.toml
